# onlyskills-repo/systemd/monster-zone-7.service
[Unit]
Description=Monster DAO Zone 7 (Prime 19)
After=network.target

[Service]
Type=simple
User=monster_zone_7
Group=monster_zone_7
SELinuxContext=system_u:system_r:monster_zone_7_t:s7

# Resource limits (DAO-governed)
CPUShares=1900
MemoryLimit=19456M
TasksMax=190
LimitNOFILE=1900

# Network namespace
PrivateNetwork=yes
IPAddressDeny=any
IPAddressAllow=10.71.7.0/24

# Security
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

ExecStart=/usr/local/bin/monster_zone_7
Restart=always

[Install]
WantedBy=multi-user.target
