# onlyskills-repo/systemd/monster-zone-1.service
[Unit]
Description=Monster DAO Zone 1 (Prime 3)
After=network.target

[Service]
Type=simple
User=monster_zone_1
Group=monster_zone_1
SELinuxContext=system_u:system_r:monster_zone_1_t:s1

# Resource limits (DAO-governed)
CPUShares=300
MemoryLimit=3072M
TasksMax=30
LimitNOFILE=300

# Network namespace
PrivateNetwork=yes
IPAddressDeny=any
IPAddressAllow=10.71.1.0/24

# Security
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

ExecStart=/usr/local/bin/monster_zone_1
Restart=always

[Install]
WantedBy=multi-user.target
